# workspace/crontab/rotordock
* * * * * rotordock php /var/www/app/cron.php 2>&1
